# simon/naemon/systemd/proxy-to-livestatus.service
[Unit]
Requires=naemon.service
After=naemon.service

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd /run/naemon/live
