# simon/naemon/systemd/status-dat-server@.service
[Service]
ExecStart=/usr/bin/cat /run/naemon/status.dat
StandardInput=socket
