# tce-operador/files/systemd/x11vnc-0@.service
[Unit]
Description=X11VNC for display :0.0

[Service]
ExecStart=/usr/bin/x11vnc -display :0.0 -inetd -quiet
StandardInput=socket
