# tce-operador/files/systemd/x11vnc-3@.service
[Unit]
Description=X11VNC for display :0.3

[Service]
ExecStart=/usr/bin/x11vnc -display :0.3 -inetd -quiet
StandardInput=socket
