# tce-operador/files/systemd/x11vnc-0.socket
[Unit]
Description=X11VNC Socket for display :0.0

[Socket]
ListenStream=5900
Accept=yes

[Install]
WantedBy=sockets.target
