# tce-operador/files/systemd/x11vnc-2.socket
[Unit]
Description=X11VNC Socket for display :0.2

[Socket]
ListenStream=5902
Accept=yes

[Install]
WantedBy=sockets.target
