# tce-operador/files/systemd/x11vnc-2@.service
[Unit]
Description=X11VNC for display :0.2

[Service]
ExecStart=/usr/bin/x11vnc -display :0.2 -inetd -quiet
StandardInput=socket
