# tce-operador/files/systemd/x11vnc-1.socket
[Unit]
Description=X11VNC Socket for display :0.1

[Socket]
ListenStream=5901
Accept=yes

[Install]
WantedBy=sockets.target
