# tce-operador/files/systemd/x11vnc-3.socket
[Unit]
Description=X11VNC Socket for display :0.3

[Socket]
ListenStream=5903
Accept=yes

[Install]
WantedBy=sockets.target
